# Docker/supervisord.conf
[supervisord]
nodaemon=true

[program:process-engine]
command=process-engine

[program:bpmn-studio]
command=bpmn-studio --port=9000 --host=0.0.0.0
